# amiya.service
[Unit]
Description=Amiya Desktop Environment
Documentation=https://github.com/yourusername/amiya
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/amiya
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

# Environment
Environment="RUST_LOG=info"
Environment="RUST_BACKTRACE=1"

[Install]
WantedBy=graphical-session.target
